# rootfs/etc/s6-overlay/s6-rc.d/bagger-setup/up
/etc/s6-overlay/scripts/bagger-setup.sh
